# conf/supervisor/dotjs.conf
[program:dotjs]
command = /home/jumski/dotfiles/vendor/dotjs-ubuntu/bin/djsd
